# server-configs/templates/worker.target.j2
[Unit]
Description=Django Worker
Requires={{ name }}-daphne.service
{% for n in range(workers) %}
Wants={{ name }}-worker@{{ n+1 }}.service{% endfor %}

[Install]
WantedBy=multi-user.target